# conf/unattended-arch-upgrade.timer
[Unit]
Description=Unattended Arch Upgrade.

[Timer]
OnCalendar=Tue,Sat *-*-* 16:30:02
Persistent=false

[Install]
WantedBy=multi-user.target
