# conf/unattended-arch-upgrade.service
[Unit]
Description=Unattended Arch Upgrade
Requires=network-online.target 

[Service]
Type=simple
ExecStart=/usr/local/bin/uau
TimeoutStopSec=900
KillMode=process
KillSignal=SIGINT
User=archupdater

